# etc/panda-webcontrol.list
# Directories
d share/www
d bin
d etc/init.d
d etc/rc.d
d lib/python/site-packages
d lib/python/site-packages/malcolm/modules/web/www
d etc/www
d share/panda-webcontrol/templates

# Python libs
B lib/python/site-packages malcolm
B lib/python/site-packages annotypes
l lib/python/site-packages/malcolm/modules/web/www/fpga_docs /opt/share/www/panda-fpga

# Startup
t bin src/panda-webcontrol.py
t etc/init.d etc/panda-webcontrol
l etc/rc.d/S400panda-webcontrol ../init.d/panda-webcontrol
l etc/rc.d/K100panda-webcontrol ../init.d/panda-webcontrol

# Templates
t etc/www etc/panda-webcontrol.index.html
t etc/www etc/panda-webcontrol.nav.html
B share/panda-webcontrol templates

# Docs
t etc/www etc/panda-webcontrol.docs.html
B share/www/panda-webcontrol malcolmjs_docs/html
